# zabbix-iostat-mon.cron
* * * * * root /usr/bin/zabbix-iostat-mon -l /var/log/zabbix-iostat-mon/zabbix-iostat-mon.log
